# cldTVUpdate.timer
[Unit]
Description=cldTV Update timer

[Timer]
OnBootSec=120min
OnUnitActivateSec=120min

[Install]
WantedBy=timers.target
