# cldTVUpdate.timer
[Unit]
Description=cldTV Update timer

[Timer]
OnBootSec=5min
OnUnitActivateSec=2h

[Install]
WantedBy=timers.target
